# gui/supervisord.conf
; supervisord config for the prusaslicer-docker :gui image.
;
; Orchestrates: Xvfb (software X server) → fluxbox (window manager) →
; x11vnc (VNC server reading the Xvfb framebuffer) → websockify (noVNC
; web frontend) → prusa-slicer (the app itself).
;
; `startsecs` is staggered so each program waits for its dependency to
; settle before launching. autorestart=true so a crash in any single
; piece doesn't take down the container.

[supervisord]
nodaemon=true
loglevel=info
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1920x1080x24 -nolisten tcp
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=":0",HOME="/root"
autorestart=true
priority=20
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:x11vnc]
# IMPORTANT: no -bg here. -bg forks x11vnc into the background and the
# original process exits — supervisord interprets that as a crash and
# eventually gives up with FATAL status. The forked daemon keeps running
# orphaned, but a fresh supervisord restart can't reach it and ends up
# fighting for port 5900 ("could not obtain listening port" in the log).
# Keep x11vnc in the foreground so supervisord manages it cleanly.
command=/usr/bin/x11vnc -display :0 -nopw -forever -shared -rfbport 5900 -o /var/log/x11vnc.log -quiet
autorestart=true
priority=30
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:novnc]
; websockify with --web=/usr/share/novnc serves the noVNC HTML at /vnc.html
; on port 8080, proxying browser WebSocket connections to localhost:5900.
command=/usr/bin/websockify --web=/usr/share/novnc 8080 localhost:5900
autorestart=true
priority=40
startsecs=4
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:prusa-slicer]
command=/usr/local/bin/prusa-slicer
environment=DISPLAY=":0",HOME="/root"
autorestart=true
priority=50
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
